# supervisor.confd/ssmtpinit.conf
[program:ssmtpinit]
command=/bin/bash -c "source /usr/local/bin/scripts/ssmtpinit"
autorestart=false
startsecs=0
priority=200